# build/susi/DEBIAN/postinst
#!/bin/bash

chmod +x /etc/init.d/susi
update-rc.d susi defaults

exit 0
